# infra/pinchtab.service
[Unit]
Description=PinchTab browser automation daemon
After=network.target

[Service]
Type=simple
User=ben
WorkingDirectory=/home/ben

# PINCHTAB_BINARY_PATH tells the pinchtab wrapper which versioned binary to use.
# Update the version string if you upgrade PinchTab.
Environment=PINCHTAB_BINARY_PATH=/home/ben/.pinchtab/bin/0.7.6/pinchtab-linux-amd64
Environment=BRIDGE_STEALTH=full

ExecStart=/home/ben/.npm-global/bin/pinchtab

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
