# deploy/xdm.service
[Unit]
Description=XDM - X/Twitter DM Auto-Reply Bot
After=network.target
# Prevent infinite rapid restart loops (max 10 restarts per 10 minutes)
StartLimitIntervalSec=600
StartLimitBurst=10

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_DIR__
ExecStart=/usr/bin/xvfb-run --auto-servernum --server-args="-screen 0 1280x1024x24" __APP_DIR__/venv/bin/python src/main.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xdm

# Clean shutdown: SIGTERM to entire cgroup (Python + Chrome subtree), then
# SIGKILL the whole cgroup after timeout. This is the default mode but stated
# explicitly — `mixed` only signals the main PID and leaves Chrome children
# orphaned in the cgroup, which causes "left-over process" warnings on restart.
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=30

# Resource limits to prevent Chrome from consuming all VPS memory
MemoryMax=2G
CPUQuota=80%
LimitNOFILE=4096

# Environment variables (uncomment and fill in to use env vars instead of config.json)
# Environment=X_USERNAME=your_username
# Environment=X_PASSWORD=your_password
# Environment=X_PASSCODE=1234
# Environment=OPENROUTER_API_KEY=sk-or-v1-xxx
# Environment=SUPABASE_URL=https://xxx.supabase.co
# Environment=SUPABASE__SECRET_KEY=xxx
# Environment=SUPABASE_DB_URL=postgresql://xxx
# Environment=TELEGRAM_BOT_TOKEN=xxx
# Environment=TELEGRAM_CHAT_ID=xxx
# Environment=FAL_KEY=xxx

[Install]
WantedBy=multi-user.target
